# deployment/recycling-api.service
[Unit]
Description=Recycling Management System - Backend API
Documentation=https://github.com/021650641/recycling-management-system
After=network.target postgresql.service
Requires=postgresql.service
StartLimitBurst=10
StartLimitIntervalSec=120

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/recycling/backend
ExecStartPre=+/bin/sh -c 'fuser -k 5000/tcp 2>/dev/null; sleep 2; exit 0'
ExecStart=/usr/bin/node dist/server.js
Restart=on-failure
RestartSec=5

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/var/www/recycling/backend/.env

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

# Resource limits
MemoryMax=512M
TasksMax=256

# Logging (use journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=recycling-api

[Install]
WantedBy=multi-user.target
